# octowatch.service
[Unit]
Description=Flash on octoprint ready

[Service]
Type=simple
ExecStart=/bin/bash /usr/bin/octowatch.sh

[Install]
WantedBy=multi-user.target
